# telegram-bot.service
[Unit]
Description=Telegram Bot Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/telegram-bot
Environment=PATH=/opt/telegram-bot/venv/bin
ExecStart=/opt/telegram-bot/venv/bin/python /opt/telegram-bot/telegram_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target